# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:sync]
command=/ocfdocs/main.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/sync/error.log
stdout_logfile=/var/log/sync/access.log