# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:sync]
command=python3 /ocfdocs/main.sh
autostart=true
autorestart=true